# Actividad4/actividad-4.service
[Unit]
Description=Servicio de prueba Actividad 4

[Service]
ExecStart=/usr/local/bin/ServicioActividad4.sh

[Install]
WantedBy=multi-user.target
